# root/etc/services.d/fava/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/bin/fava -H 0.0.0.0
